# pianoteq-autolink/.config/systemd/user/pianoteq-autolink.service
[Unit]
Description=Auto-connect Pianoteq to iConnectAUDIO4+
After=pipewire.service wireplumber.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pianoteq-autolink.sh
